# profile.d/nvm.sh
#!/bin/sh

# 配置 nvm 自启 和 镜像
# export NVM_DIR="/home/ubuntu/.nvm"
# [ -s "$NVM_DIR/nvm.sh" ] && \. "$NVM_DIR/nvm.sh"  # This loads nvm
# [ -s "$NVM_DIR/bash_completion" ] && \. "$NVM_DIR/bash_completion"  # This loads nvm bash_completion
# export NVM_NODEJS_ORG_MIRROR=https://npmmirror.com/mirrors/node/ # This set the customer mirror for node
